# doc/examples/croncho.cron
*/1  * * * * /usr/local/icn/croncho/bin/croncho.pl
